# docker-files/supervisord/_program-mqtt.conf
[program:mqtt]
command=sh -c "(php artisan -vvv app:mqtt-loop) || sleep 30"
autorestart=true
